# setup-hawkradius-host/files/systemd/apt-daily-upgrade.service.d/override.conf
[Unit]
Wants=ghost-update.service
Before=ghost-update.service
